# deploy/gunicorn.service
[Unit]
Description=Gunicorn Sport Plus Management
After=network.target postgresql.service

[Service]
User=spm
Group=www-data
WorkingDirectory=/srv/spm
EnvironmentFile=/etc/spm.env
RuntimeDirectory=spm
RuntimeDirectoryMode=0755
ExecStart=/srv/spm/.venv/bin/gunicorn --workers 3 --bind unix:/run/spm/gunicorn.sock --access-logfile - --error-logfile - spm.wsgi:application
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
